# 85-mymouse.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/lib/udev/mymouse.sh"
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/lib/udev/mymouse.sh"
